# systemd/portal.service
[Unit]
Description=DB Studio Customer Portal (main app)
After=network.target postgresql.service portal-pdf.service
Requires=postgresql.service
PartOf=portal-pdf.service

[Service]
Type=simple
User=portal-app
Group=portal-app
WorkingDirectory=/opt/dbstudio_portal
ExecStart=/opt/dbstudio_portal/.node/bin/node server.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/opt/dbstudio_portal/.env

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ReadWritePaths=/var/lib/portal /var/log/portal /run/portal-pdf
ReadOnlyPaths=/opt/dbstudio_portal
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=true
LockPersonality=true
# MemoryDenyWriteExecute=true — DELIBERATELY OMITTED. Node 20 V8 JIT requires WX pages,
# same constraint that already forced its omission on portal-pdf.service.
RestrictRealtime=true
SystemCallArchitectures=native
# SystemCallFilter — DELIBERATELY OMITTED. Even @system-service crashes Node 20 (status=31/SYS).
# The other hardening (NoNewPrivileges, ProtectSystem=strict, ProtectHome, PrivateTmp,
# RestrictAddressFamilies, ProtectKernel*, RestrictNamespaces, LockPersonality, RestrictRealtime)
# remains strong. Seccomp can be re-added later once we identify exactly which calls Node needs
# and write a custom allowlist via SystemCallFilter=… narrowed to that set.

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=portal

[Install]
WantedBy=multi-user.target
